# deploy/efl-fetch-proxy/efl-fetch-proxy.service
[Unit]
Description=IntelliWatt EFL Fetch Proxy (WAF fallback)
After=network.target

[Service]
Type=simple
User=deploy

# Repo standard on this droplet:
# - app repo lives at /home/deploy/apps/intelliwatt
# - per-service env lives at /home/deploy/.<service>.env
WorkingDirectory=/home/deploy/apps/intelliwatt/deploy/efl-fetch-proxy
ExecStart=/usr/bin/node dist/server.js
Restart=always
RestartSec=2

# Load shared IntelliWatt env if present (optional)
EnvironmentFile=-/home/deploy/.intelliwatt.env

# Keep secrets out of repo; load from a deploy-owned env file.
EnvironmentFile=/home/deploy/.efl-fetch-proxy.env

# Hardening (minimal; adjust if needed).
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

