# deploy/smt/smt-ingest.service
[Unit]
Description=IntelliWatt SMT ingest (SFTP fetch -> droplet upload server POST)
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=deploy
Group=deploy
EnvironmentFile=/etc/default/intelliwatt-smt
WorkingDirectory=/home/deploy/apps/intelliwatt
ExecStart=/bin/bash /home/deploy/apps/intelliwatt/deploy/smt/fetch_and_post.sh

# Auto-restart on failure to ensure SMT ingest always retries
# This service runs on demand and must always be available
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
