# deploy/smt/smt-ingest.timer
[Unit]
Description=Run IntelliWatt SMT ingest every 30 minutes

[Timer]
OnBootSec=2m
OnUnitActiveSec=30m
AccuracySec=1m
Persistent=true
Unit=smt-ingest.service

[Install]
WantedBy=timers.target
