# deploy/droplet/efl-pdftotext.service
[Unit]
Description=IntelliWatt EFL pdftotext Helper (Python on 8095)
After=network.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/home/deploy/apps/intelliwatt

# EFL pdftotext specific environment (do not use /home/deploy/.intelliwatt.env; it may contain `export` lines,
# which systemd will ignore and spam logs).
EnvironmentFile=-/home/deploy/.efl-pdftotext.env

Environment=EFL_PDFTEXT_PORT=8095

# Start the helper
ExecStart=/usr/bin/env python3 /home/deploy/apps/intelliwatt/deploy/droplet/efl_pdftotext_server.py

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
